# systemd/system/snap-google\x2dcloud\x2dsdk-132.mount
[Unit]
Description=Mount unit for google-cloud-sdk, revision 132
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/google-cloud-sdk_132.snap
Where=/snap/google-cloud-sdk/132
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
